# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm

[program:npm]
command=npm run dev -- --host 0.0.0.0
directory=/var/www
